# config/includes.chroot_before_packages/etc/skel/.config/autostart/ufstartupcheck.desktop
[Desktop Entry]
Exec=/usr/share/ufhelper/uflinux-helper-startupcheck
Icon=/usr/share/ufhelper/ufhelperlogo.svg
Name=uflinuxstartupchecker
StartupNotify=false
Terminal=false
Type=Application
