# config/includes.chroot_before_packages/etc/cron.d/uflinuxofflinepackages
@reboot root /var/uflinuxofflinepackages/sources.sh
