# config/includes.chroot_before_packages/etc/cron.d/uflinuxofflinepackages
@reboot root /bin/sh /var/uflinuxofflinepackages/sources.sh
