# config/includes.chroot_before_packages/etc/cron.d/uflinuxofflinepackages
@reboot root chmod +x /var/uflinuxofflinepackages/sources.sh && /var/uflinuxofflinepackages/sources.sh
